# systemd/aco-panel.service
[Unit]
Description=ACO Maintenance Panel Web Server
Documentation=https://github.com/xofyy/aco-panel
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aco-panel

# Environment
Environment="PATH=/opt/aco-panel/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Start Flask with Gunicorn
ExecStart=/opt/aco-panel/venv/bin/python -m gunicorn \
    --bind 0.0.0.0:4444 \
    --workers 2 \
    --threads 4 \
    --timeout 120 \
    --access-logfile /var/log/aco-panel-access.log \
    --error-logfile /var/log/aco-panel-error.log \
    --capture-output \
    app.main:app

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
