# systemd/kiosk-panel.service
[Unit]
Description=Kiosk Setup Panel Web Server
Documentation=https://github.com/xofyy/kiosk-setup-panel
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/kiosk-setup-panel

# Environment
Environment="PATH=/opt/kiosk-setup-panel/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Gunicorn ile Flask başlat
ExecStart=/opt/kiosk-setup-panel/venv/bin/python -m gunicorn \
    --bind 0.0.0.0:8080 \
    --workers 2 \
    --threads 4 \
    --timeout 120 \
    --access-logfile /var/log/kiosk-panel-access.log \
    --error-logfile /var/log/kiosk-panel-error.log \
    --capture-output \
    app.main:app

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
